# root.crontab
* * * * * /home/watchdoguser/evil-process-watchdog-killer/killevils.sh > /dev/null
