# watchdoguser.crontab
* * * * *  ~/evil-process-watchdog-killer/watchdog.sh
* * * * *  ~/evil-process-watchdog-killer/top1.sh
